# config/supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/bin/bash -c "/usr/local/bin/docker-entrypoint.sh mongod --keyFile /app/mongo/keyfile --bind_ip_all --replSet rs0"
priority=1
startsecs=30
redirect_stderr=false
stdout_logfile=/app/logs/mongod.log

[program:init-mongodb]
command=/bin/bash -c "mongo -u $MONGO_INITDB_ROOT_USERNAME -p $MONGO_INITDB_ROOT_PASSWORD < /app/mongo/mongo-init.js"
priority=2
redirect_stderr=false
stdout_logfile=/app/logs/mongod.log

[program:4km-service]
command=/bin/bash -c "/app/4km-service.sh"
priority=10
startsecs=10
redirect_stderr=false
stdout_logfile=/app/logs/4km-service.log

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;' -c /etc/nginx/nginx.conf
priority=20
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/nginx.log
stderr_logfile=/app/logs/nginx.err.log
